# vault/supervisor/vault.conf
[program:vault]
command=vault server -config=/root/vault-config.hcl
autostart=true
autorestart=true
stderr_logfile=/var/log/vault.out.err
stdout_logfile=/var/log/vault.out.log
